# scripts/opencode-server.service
[Unit]
Description=OpenCode Headless HTTP Server
Documentation=https://opencode.ai
# opencode binds to 127.0.0.1 only — no network-level dependency needed.
# After=network.target is a system target and ignored in user units anyway.

[Service]
Type=simple
# No User= — user units run as the owning user automatically.
# %h expands to the user's home directory at activation time.
WorkingDirectory=%h

# Port defaults to 4096. Override without editing this file:
#   echo "OPENCODE_PORT=4097" > ~/.config/opencode/opencode-server.env
#   systemctl --user restart opencode-server
Environment=OPENCODE_PORT=4096
# EnvironmentFile: '-' prefix = silently ignore if file doesn't exist.
EnvironmentFile=-%h/.config/opencode/opencode-server.env

# Full path required — opencode is installed per-user under ~/.opencode/bin.
ExecStart=%h/.opencode/bin/opencode serve \
    --port ${OPENCODE_PORT} \
    --print-logs \
    --log-level INFO

# Restart the service if it exits for any reason other than a clean stop.
Restart=on-failure
RestartSec=5s

# Cap restart storms: at most 3 restarts within 60 seconds, then give up
# and leave the service in a failed state for manual inspection.
StartLimitIntervalSec=60
StartLimitBurst=3

# Writes to the per-user journal (volatile by default — cleared on reboot).
# Query: journalctl --user -u opencode-server [-f]
#        journalctl --user -u opencode-server --since "1 hour ago"
StandardOutput=journal
StandardError=journal
SyslogIdentifier=opencode-server
# HOME is set automatically by the user manager — no Environment=HOME= needed.

[Install]
# default.target is the correct anchor for user services (not multi-user.target).
# Service starts at login; survives logout after running once: loginctl enable-linger
WantedBy=default.target
